# evolution/deploy/evolution-runner.service
# Systemd unit for the evolution batch runner (Type=oneshot, triggered by timer).
# Runs the batch runner script to claim and execute pending evolution runs.

[Unit]
Description=Evolution Pipeline Batch Runner
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/explainanything
# No EnvironmentFile — script loads .env.local and .env.evolution-prod directly via dotenv
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/npx tsx evolution/scripts/processRunQueue.ts --max-runs 10 --parallel 10
TimeoutStartSec=7200
User=evolution
Group=evolution

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=evolution-runner
